# files/etc/cron.d/puppet
*/5 * * * * root /etc/puppet/files/bin/puppet-client.sh --logdest syslog > /dev/null 2>&1
